# systemd/cwspots-kiwis.timer
[Unit]
Description=Weekly KiwiSDR station list refresh

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
