# systemd/cwspots-skcc.timer
[Unit]
Description=Weekly SKCC member roster refresh (Friday mornings)

[Timer]
OnCalendar=Fri *-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
