# systemd/cwspots-ctydat.timer
[Unit]
Description=Monthly refresh of cty.dat DXCC data

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
