# systemd/cwspots-kiwis.service
[Unit]
Description=Fetch KiwiSDR station list
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=cwspots
WorkingDirectory=/opt/cwspots
ExecStart=/opt/cwspots/venv/bin/python /opt/cwspots/fetch_kiwis.py
