# systemd/cwspots-tlscert.timer
[Unit]
Description=Weekly Tailscale TLS cert renewal check

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
